# ci-cd.conf
[program:ci-cd]
user=user
directory=/home/user/ci-cd
command=/home/user/ci-cd/venv/bin/gunicorn --bind=:8000 --log-level debug --workers=3 ci_cd.wsgi

autostart=true
autorestart=true
stdout_logfile=/home/user/logs/gunicorn.log
stderr_logfile=/home/user/logs/gunicorn.err.log
